# picamera2-webstream.service
[Unit]
Description=PiCamera2 Web Stream Service
After=network.target

[Service]
Type=simple
User=ian
Group=ian
WorkingDirectory=/home/ian/picamera2-webstream
Environment=PATH=/home/ian/picamera2-webstream/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ian/picamera2-webstream/venv/bin/python /home/ian/picamera2-webstream/examples/stream_service.py
Restart=always
RestartSec=3
StandardOutput=append:/var/log/picamera2-webstream.log
StandardError=append:/var/log/picamera2-webstream.log

[Install]
WantedBy=multi-user.target
